# lesson99-DIPLOM/vm/ansible/roles/15_dataservers/files/borg-backup.timer
[Unit]
Description=Automated Borg Backup Timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target